# config/supervisor-local/lsports.conf
[group:LSPORTS]
programs=loaderprematch,loaderlive,stmfixture,stmbet_01,stmbet_02,stmbet_03,stmbet_04,stmbet_05,stmbet_06,stmbet_07,stmbet_08,stmbet_09,stmbet_10

[program:loaderprematch]
command= nice -n 22 /usr/bin/php /home/alexis/projects/perso/admin_affiliation/bin/console  app:testlssportprematch
user=www-data
redirect_stderr=false
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
priority=1
autostart=false

[program:loaderlive]
command= nice -n 22 /usr/bin/php /home/alexis/projects/perso/admin_affiliation/bin/console  app:testlssportinplay
user=www-data
redirect_stderr=false
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
priority=1
autostart=false

[program:stmfixture]
command= nice -n 22 /usr/bin/php /home/alexis/projects/perso/admin_affiliation/bin/console  rabbitmq:consumer lsports_fixture -m 100
user=www-data
redirect_stderr=false
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
priority=1
autostart=false

[program:stmbet_01]
command= nice -n 22 /usr/bin/php /home/alexis/projects/perso/admin_affiliation/bin/console  rabbitmq:consumer lsports_bet -m 100
user=www-data
redirect_stderr=false
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
priority=1
autostart=false

[program:stmbet_02]
command= nice -n 22 /usr/bin/php /home/alexis/projects/perso/admin_affiliation/bin/console  rabbitmq:consumer lsports_bet -m 100
user=www-data
redirect_stderr=false
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
priority=1
autostart=false

[program:stmbet_03]
command= nice -n 22 /usr/bin/php /home/alexis/projects/perso/admin_affiliation/bin/console  rabbitmq:consumer lsports_bet -m 100
user=www-data
redirect_stderr=false
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
priority=1
autostart=false

[program:stmbet_04]
command= nice -n 22 /usr/bin/php /home/alexis/projects/perso/admin_affiliation/bin/console  rabbitmq:consumer lsports_bet -m 100
user=www-data
redirect_stderr=false
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
priority=1
autostart=false

[program:stmbet_05]
command= nice -n 22 /usr/bin/php /home/alexis/projects/perso/admin_affiliation/bin/console  rabbitmq:consumer lsports_bet -m 100
user=www-data
redirect_stderr=false
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
priority=1
autostart=false

[program:stmbet_06]
command= nice -n 22 /usr/bin/php /home/alexis/projects/perso/admin_affiliation/bin/console  rabbitmq:consumer lsports_bet -m 100
user=www-data
redirect_stderr=false
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
priority=1
autostart=false

[program:stmbet_07]
command= nice -n 22 /usr/bin/php /home/alexis/projects/perso/admin_affiliation/bin/console  rabbitmq:consumer lsports_bet -m 100
user=www-data
redirect_stderr=false
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
priority=1
autostart=false

[program:stmbet_08]
command= nice -n 22 /usr/bin/php /home/alexis/projects/perso/admin_affiliation/bin/console  rabbitmq:consumer lsports_bet -m 100
user=www-data
redirect_stderr=false
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
priority=1
autostart=false

[program:stmbet_09]
command= nice -n 22 /usr/bin/php /home/alexis/projects/perso/admin_affiliation/bin/console  rabbitmq:consumer lsports_bet -m 100
user=www-data
redirect_stderr=false
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
priority=1
autostart=false

[program:stmbet_10]
command= nice -n 22 /usr/bin/php /home/alexis/projects/perso/admin_affiliation/bin/console  rabbitmq:consumer lsports_bet -m 100
user=www-data
redirect_stderr=false
autorestart=true
startsecs=2
startretries=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
priority=1
autostart=false